# openfortivpn/ip_report.cron
# report my IP to sas-api every minute
* * * * * root /home8/lcrown/dotfiles/openfortivpn/report.py
